# deploy/claude-rpg.service
[Unit]
Description=Claude RPG - Claude Code Companion
After=network.target

[Service]
Type=simple
WorkingDirectory=%h/dev/whoabuddy/claude-rpg/server-v2
Environment=NODE_ENV=production
Environment=CLAUDE_RPG_PORT=4011
Environment=PATH=%h/.bun/bin:/usr/local/bin:/usr/bin:/bin
ExecStart=%h/.bun/bin/bun run index.ts
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=default.target
